# services/disable-usb-wakeup.service
[Unit]
Description=Disable external USB device wakeup

[Service]
Type=oneshot
ExecStart=/bin/sh -c "echo XHCI > /proc/acpi/wakeup"

[Install]
WantedBy=multi-user.target
